# petabyte_out/lumaris_agent/petabyte-agent.service
[Unit]
Description=Petabyte GPU node agent
After=network.target docker.service
Wants=docker.service

[Service]
WorkingDirectory=/opt/petabyte-agent
EnvironmentFile=/etc/petabyte/agent.env
ExecStart=/opt/petabyte-agent/.venv/bin/python main.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
